# rc.local
#!/bin/sh

sleep 25
/etc/init.d/passwall restart

exit 0
