# files/etc/rc.local
#!/bin/sh
sleep 20
/etc/init.d/passwall restart
exit 0
